# config/supervisor/conf.d/gofaxd.conf
[program:gofaxd]
user=uucp
command=/usr/local/sbin/gofaxd
autorestart=true
